# code/deploy/launchd/com.largefilebridge.pin.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- The 15-minute background pin/add-to-IPFS pass (storage_local.mdx §13). Installed by the web app / justfile. -->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.largefilebridge.pin</string>
  <key>ProgramArguments</key>
  <array>
    <string>__NODE_BIN__</string>
    <string>__TRIGGER_SCRIPT__</string>
    <string>pin</string>
    <string>__API_PORT__</string>
  </array>
  <key>StartInterval</key><integer>900</integer>
  <key>RunAtLoad</key><false/>
  <key>StandardOutPath</key><string>__LOG_OUT__</string>
  <key>StandardErrorPath</key><string>__LOG_ERR__</string>
</dict>
</plist>
